# muondetector-daemon/dpkg-deb-packaging/raspbian/muondetector-daemon/DEBIAN/postinst
#!/bin/bash
sed -i 's/pigpiod -l/pigpiod -s 1/g' /lib/systemd/system/pigpiod.service
systemctl enable pigpiod.service
systemctl enable muondetector-daemon.service
