# muondetector-daemon/dpkg-deb-packaging/raspbian/muondetector-daemon/DEBIAN/postinst
#!/bin/bash
systemctl enable pigpiod.service
systemctl enable muondetector-daemon.service
