# muondetector-daemon/dpkg-deb-packaging/raspbian/muondetector-daemon/DEBIAN/prerm
#!/bin/bash
systemctl disable muondetector-daemon.service
systemctl disable pigpiod.service
