# daemon/config/muondetector-daemon.service
[Unit]
Description=muondetector-daemon - Daemon for custom muondetector board
<<<<<<< HEAD:daemon/config/muondetector-daemon.service
After=sockets.target network-online.target
=======
Requires=pigpiod.service
BindsTo=pigpiod.service
After=pigpiod.service sockets.target network-online.target
>>>>>>> dev:daemon/config/muondetector-daemon.service

[Service]
Type=simple
User=muonuser
ExecStart=/usr/bin/muondetector-daemon -c
ExecStop=/bin/kill $MAINPID
Restart=always
TimeoutSec=2

[Install]
WantedBy=multi-user.target
